# app-tests/sample_service/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:flask]
command=flask run --host=0.0.0.0 --port=5000